# packaging/sound-server.path
[Unit]
Description=Path activation for the "sound-server" programme

[Path]
PathExists=/tmp/.pa_ready
